# docker/s6-overlay/s6-rc.d/command-monitor/run
#!/command/with-contenv sh
exec bash -c 'cd "${APP_BASE_DIR:-/var/www/html}" && exec php tempest command:monitor'
