# docker/s6-overlay/s6-rc.d/schedule/run
#!/command/with-contenv sh
exec bash -c 'cd "${APP_BASE_DIR:-/var/www/html}" && while true; do php tempest schedule:run; sleep 60; done'
